# packaging/sentinel-rotate.timer
[Unit]
Description=Daily rotation timer for SentinelEdgeAI health credentials

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
