# packaging/sentinel-firewall-rollback.timer
[Unit]
Description=Run SentinelEdgeAI Firewall Rollback Guard every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=sentinel-firewall-rollback.service

[Install]
WantedBy=timers.target
