# packaging/sentinel-health-vault.service
[Unit]
Description=Fetch Sentinel health secret from Vault
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/fetch-health-from-vault.sh
# Default to http (Vault dev server uses HTTP). Can be overridden in environment.
Environment=VAULT_ADDR=http://127.0.0.1:8200
EnvironmentFile=/etc/sentinel/vault.env

[Install]
WantedBy=multi-user.target
